# vnpy_engine/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/vnpy/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:api-gateway]
command=python -m vnpy_engine.vnpy_local.api_gateway_enhanced
directory=/vnpy
environment=PYTHONPATH="/vnpy",REDIS_HOST="%(env_REDIS_HOST)s",REDIS_PORT="%(env_REDIS_PORT)s",TRADING_MODE="%(env_TRADING_MODE)s",BINANCE_API_KEY="%(env_BINANCE_API_KEY)s",BINANCE_SECRET_KEY="%(env_BINANCE_SECRET_KEY)s",BINANCE_SERVER="%(env_BINANCE_SERVER)s",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=15
startretries=5
exitcodes=0
stdout_logfile=/vnpy/logs/api-gateway.log
stderr_logfile=/vnpy/logs/api-gateway.err
killasgroup=true
priority=100

[program:watchdog]
command=python -m vnpy_engine.vnpy_local.watchdog
directory=/vnpy
environment=PYTHONPATH="/vnpy",REDIS_HOST="%(env_REDIS_HOST)s",REDIS_PORT="%(env_REDIS_PORT)s",TRADING_MODE="%(env_TRADING_MODE)s",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=5
exitcodes=0
stdout_logfile=/vnpy/logs/watchdog.log
stderr_logfile=/vnpy/logs/watchdog.err
killasgroup=true
priority=200

[group:vnpy-services]
programs=api-gateway,watchdog
priority=999
